# debian/postinst
#!/bin/bash -e

/usr/lib/pt-device-manager/config/xhost-setup

echo "Unmasking device manager service..."
systemctl unmask pt-device-manager

echo "Enabling device manager service..."
systemctl enable pt-device-manager

echo "Starting device manager service..."
systemctl start pt-device-manager
