# debian/prerm
#!/bin/bash -e

echo "Stopping device manager service..."
systemctl stop pt-device-manager

echo "Disabling device manager service..."
systemctl disable pt-device-manager
